# transaction-kafka-producer/crontab
0 0,12 * * * cd /app && python main.py >> /var/log/cron.log 2>&1
